# deploy/systemd/boardline-one8-daily.timer
[Unit]
Description=Run Boardline One8 Daily Sync Every Day

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
