# deploy/systemd/boardline-one8-daily.service
[Unit]
Description=Boardline One8 Daily Sync
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/boardline_spider
ExecStart=/bin/bash /root/boardline_spider/run_one8_daily.sh
